# debian/cron.daily
#!/bin/sh


test -x /usr/sbin/logrotate -a ! -d /lib/systemd/system || exit 0
/usr/sbin/logrotate /etc/logrotate.conf
